# modules/game/scripts/jobs.cron
*/5 * * * * /home/ubuntu/${game_name}/backups.sh >/home/ubuntu/${game_name}/backup.log 2>&1
*/5 * * * * /home/ubuntu/${game_name}/duck.sh >/dev/null 2>&1
@reboot /home/ubuntu/${game_name}/duck.sh >/dev/null 2>&1
